# supervisord.conf
[supervisord]
nodaemon=true

[program:network]
user=root
command=iptables -t nat -A PREROUTING -i eth0 -p tcp --dport 80 -j DNAT --to 10.13.37.100:80
autorestart=false
startsecs=0